# collector/cdn-stats.timer
[Unit]
Description=Run CDN stats collector every 10 minutes

[Timer]
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
